# automation/yara-monitoring-ansible/templates/yara_monitor.conf.upstart.j2
description "Yara Monitoring for Cowrie"
start on runlevel [2345]
stop on shutdown

script
    sudo /opt/yara/yara_monitor_linux.sh
end script
