# automation/cowrie-ansible/templates/cowrie.service.socket.j2
# This cowrie.socket file has been automatically generated by Ansible

[Unit]
Description=Cowrie SSH honeypot socket
[Socket]
ListenStream=0.0.0.0:{{ cowrie_ssh_port }}
[Unit]
Description=Cowrie Telnet honeypot socket
[Socket]
ListenStream=0.0.0.0:{{ cowrie_telnet_port }}
[Install]
WantedBy=sockets.target